# deploy/thesis-api.service
[Unit]
Description=Thesis API
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/MDsystem/apps/api
EnvironmentFile=-/root/MDsystem/apps/api/.env
ExecStart=/usr/bin/node /root/MDsystem/apps/api/dist/server.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
